# build/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
user=root  # Run Supervisor as a non-root user

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:app]
command=python app.py  # Adjust if your app has a different entry point
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
user=nobody  # Run the app as a non-root user

[program:monitor]
command=python monitor.py
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
user=nobody  # Run the monitor as a non-root user

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
